# archlinux/systemd/system/kbd_idle.service
[Unit]
Description=Disables keyboard backlight when inactive.

[Service]
ExecStart=/home/ak/scripts/kbd_idle

[Install]
WantedBy=multi-user.target
