# archlinux/systemd/system/kill-wan.service
[Unit]
Description=Kill WAN Connection
Before=sleep.target

[Service]
Type=oneshot
ExecStart=-/sbin/rmmod xmm7360

[Install]
WantedBy=sleep.target